# deploy/supervisor.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:yukti]
command=uv run python -m yukti --mode %(ENV_MODE)s
directory=/app
autostart=true
autorestart=true
startretries=5
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/yukti.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/yukti_err.log
stderr_logfile_maxbytes=10MB
environment=
    PYTHONUNBUFFERED="1",
    PYTHONDONTWRITEBYTECODE="1"

[eventlistener:crash_alert]
command=python3 /app/scripts/crash_alert.py
events=PROCESS_STATE_FATAL
